# meta-mcst/meta-common/recipes-mcst/allow-host-boot/files/allow_host_boot.service
[Unit]
Description=Allow Host Boot
After=mcst-fruid.service

[Service]
ExecStart=/libexec/allow_host_boot
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=basic.target
