# meta-depo/meta-dacn469535003v2/recipes-depo/power-manager/files/host-poweroff.service
[Unit]
Description=Stop Host
Conflicts=obmc-host-startmin@0.target

[Service]
Type=oneshot
ExecStart=/usr/bin/host off
SyslogIdentifier=buttonscripts
