# meta-depo/meta-dacn469535003v2/recipes-depo/power-manager/files/host-reset.service
[Unit]
Description=Reset Host
Conflicts=obmc-chassis-poweroff@0.target host-poweroff-hard.service host-poweroff.service host-poweron.service 

[Service]
Type=oneshot
ExecStart=/usr/bin/host reset
SyslogIdentifier=buttonscripts
