# meta-depo/meta-dacn469535003v2/recipes-depo/power-manager/files/host-poweroff-hard.service
[Unit]
Description=Unconditionally Stop Host
Conflicts=obmc-host-startmin@0.target

[Service]
Type=oneshot
ExecStart=/usr/bin/host off hard
SyslogIdentifier=buttonscripts
