# meta-mcst/meta-common/recipes-mcst/tmux-logger/files/logrotate-sol.timer
[Unit]
Description=Rotate SoL log each minute

[Timer]
Unit=logrotate-sol
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
