# meta-depo/meta-dazn/recipes-depo/power-manager/files/host-poweroff-hard.service
[Unit]
Description=Unconditionally Stop Host
Conflicts=host-poweroff.service host-poweron.service host-reset.service

[Service]
Type=oneshot
ExecStart=/usr/bin/env host off hard
SyslogIdentifier=buttonscripts
